# files/cron-supervisor.ini
[program:cron]
priority=10
command=crond -n
user=root
stopwaitsecs = 3600
stopasgroup = false
killasgroup = true
stopsignal=QUIT
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
user=root
