# files/openbox-supervisor.ini
[program:openbox]
priority=25
command=/usr/bin/openbox-session
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
user=developer
environment=DISPLAY=":1",HOME="/home/developer"
autostart=true
autorestart=true
stopsignal=QUIT
