# files/ssh-supervisor.ini
[program:sshd]
priority=15
command=/usr/bin/sshd -D
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
user=root
autostart=true
autorestart=true
stopsignal=QUIT
