# files/x11vnc-supervisor.ini
[program:x11vnc]
priority=30
directory=/home/developer
environment=HOME="/home/developer",USER="developer"
command=x11vnc -forever -shared -rfbport 5900 -cursor arrow -arrow 3 -ncache 10 -noipv6 -no6 -display :1
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
user=developer
autostart=true
autorestart=false
stopsignal=QUIT
