# files/xvfb-supervisor.ini
[program:xvfb]
priority=20
command=/usr/bin/Xvfb :1 -screen 0 1152x864x16
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
user=root
autostart=true
autorestart=true
stopsignal=QUIT
